# app/cron/eramba-hourly.service
[Unit]
Description=Eramaba hourly job runner

[Service]
Type=oneshot
ExecStart=/usr/bin/sh /eramba_community/cron/eramba_hourly.sh

[Install]
WantedBy=multi-user.target
